# systemd/machinevisionflow-nodered.service
[Unit]
Description=Machine Vision Flow - Node-RED
Documentation=https://github.com/your-repo/MachineVisionFlow
After=network.target machinevisionflow-backend.service
Requires=machinevisionflow-backend.service

[Service]
Type=forking
User=cnc
Group=cnc
WorkingDirectory=/home/cnc/.node-red

# Runtime directories (shared with backend service)
RuntimeDirectory=machinevisionflow
LogsDirectory=machinevisionflow

# Environment variables for production paths
Environment="LOG_DIR=/var/log/machinevisionflow"
Environment="RUN_DIR=/run/machinevisionflow"
Environment="NODE_RED_HOME=/home/cnc/.node-red"

# Start command
ExecStart=/usr/bin/node-red
ExecStop=/bin/kill -TERM $MAINPID

# PID file
PIDFile=/run/machinevisionflow/node-red.pid

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Restart policy
Restart=on-failure
RestartSec=10s

# Logging
StandardOutput=append:/var/log/machinevisionflow/node-red.log
StandardError=append:/var/log/machinevisionflow/node-red.log

[Install]
WantedBy=multi-user.target
